# tools/depends/target/nghttp2/iphoneos16.1_arm64-target-debug/build/contrib/nghttpx-upstart.conf
# vim: ft=upstart:

description "HTTP/2 reverse proxy"

start on runlevel [2]
stop on runlevel [016]

exec /Users/Shared/xbmc-depends/iphoneos16.1_arm64-target-debug/bin/nghttpx
