# tools/depends/target/samba-gplv3/iphoneos16.1_arm64-target-debug/bin/default/packaging/systemd/nmb.service
[Unit]
Description=Samba NMB Daemon
Documentation=man:nmbd(8) man:samba(7) man:smb.conf(5)
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=notify
PIDFile=/Users/Shared/xbmc-depends/iphoneos16.1_arm64-target-debug/var/run/samba/nmbd.pid
EnvironmentFile=-/Users/Shared/xbmc-depends/iphoneos16.1_arm64-target-debug/etc/sysconfig/samba
ExecStart=/Users/Shared/xbmc-depends/iphoneos16.1_arm64-target-debug/sbin/nmbd --foreground --no-process-group $NMBDOPTIONS
ExecReload=/bin/kill -HUP $MAINPID
LimitCORE=infinity


[Install]
WantedBy=multi-user.target
